# local/oled-mon/oled_mon_delay.rpi-1.local.service
[Unit]
Description=oled_mon_delay

[Service]
Type=oneshot
ExecStart=/bin/sleep 20

[Install]
WantedBy=default.target
